# jeedom.cron
* * * * * http /usr/bin/php /var/http/jeedom/core/php/jeeCron.php > /dev/null
