# jeedom.cron
* * * * * http nice -n 19 /usr/bin/php /srv/http/jeedom/core/php/jeeCron.php > /dev/null 2>&1
